# bin/crontab
# run update script every minute
 * *  * * * cd /etc/libsites && . ~/.bashrc && ./bin/update
# get zdb data and create dump every day
53 09 * * * cd /etc/libsites && . ~/.bashrc && make zdb dump >> update.log
